# supervisord.upstream.conf
###############################################################################
# REQ-17 — supervisord config for upstream SAJHA + agent + nginx
#
# 1. sajha-upstream (priority 1) — FastAPI/uvicorn on 127.0.0.1:3002
#    Started with our tools-pack overlay via env vars.
# 2. agent (priority 2) — FastAPI/uvicorn on 127.0.0.1:8000, talks to upstream via JWT.
# 3. nginx (priority 3) — reverse proxy on $PORT.
###############################################################################

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

# ── SAJHA Upstream (priority 1) ─────────────────────────────────────────────
[program:sajha-upstream]
command=python /app/sajhamcpserver-upstream/run_server.py
directory=/app/sajhamcpserver-upstream
environment=SAJHA_SERVER_HOST="127.0.0.1",SAJHA_SERVER_PORT="3002",SAJHA_CONFIG_TOOLS_DIR="/app/tools-pack/configs",PYTHONPATH="/app/tools-pack:/app",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=8
startretries=3
priority=1
stdout_logfile=/var/log/sajha-upstream.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/sajha-upstream_err.log
stderr_logfile_maxbytes=5MB

# ── Agent Server (priority 2) ───────────────────────────────────────────────
[program:agent]
command=uvicorn agent_server:app --host 127.0.0.1 --port 8000 --workers 1
directory=/app
environment=SAJHA_BASE_URL="http://127.0.0.1:3002",SAJHA_AUTH_MODE="jwt",SAJHA_ADMIN_USER="admin",SAJHA_ADMIN_PASS="admin123",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=2
stdout_logfile=/var/log/agent.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/agent_err.log
stderr_logfile_maxbytes=5MB

# ── nginx (priority 3) ──────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
directory=/app
autostart=true
autorestart=true
startsecs=2
startretries=3
priority=3
stdout_logfile=/var/log/nginx_access.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/nginx_err.log
stderr_logfile_maxbytes=5MB
